# webhook/ramboq_dev.service
[Unit]
Description=RamboQuant API + SvelteKit SPA (dev — dev.ramboq.com)
After=network.target redis.service
Wants=redis.service

[Service]
WorkingDirectory=/opt/ramboq_dev
User=www-data
Group=www-data
ExecStart=/bin/bash -c 'source /opt/ramboq_dev/venv/bin/activate && \
  uvicorn api.app:app \
    --host 0.0.0.0 \
    --port 8001 \
    --workers 1 \
    --log-level info \
    2>&1 | tee -a /opt/ramboq_dev/.log/error_file'
Restart=always
RestartSec=5
Environment=REDIS_URL=redis://localhost:6379

[Install]
WantedBy=multi-user.target
